# crontab
52 * * * * sudo /home/pi/working/unifi-api/check-ip.sh
